# netconf/waco-net-conf.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/waco-net-conf.sh

[Install]
WantedBy=default.target
